# root/etc/uci-defaults/luci-qemu
#!/bin/sh

# Create necessary directories
mkdir -p /storage/qemu

# Make init script executable
chmod +x /etc/init.d/qemu

exit 0